# docker_config/supervisor.conf
[program:django]
command=/code/docker_config/start.sh
autostart=true
autorestart=true
[program:celery]
command=/code/docker_config/start_cel.sh
autostart=true
autorestart=true
